# blockChina.iptables
### BEGIN INIT INFO
# Provides:          scriptname
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

iptables -A INPUT -s 182.74.167.32/27 -j DROP
iptables -A INPUT -s 183.0.0.0/10 -j DROP
iptables -A INPUT -s 221.224.0.0/13 -j DROP
iptables -A INPUT -s 222.184.0.0/13 -j DROP
iptables -A INPUT -s 58.208.0.0/12 -j DROP
iptables -A INPUT -s 86.105.195.0/24 -j DROP
iptables -A INPUT -s 112.0.0.0/8 -j DROP
iptables -A INPUT -s 59.40.0.0/13 -j DROP
iptables -A INPUT -s 121.16.0.0/13 -j DROP
iptables -A INPUT -s 106.2.0.0/15 -j DROP
iptables -A INPUT -s 116.16.0.0/12 -j DROP
iptables -A INPUT -s 123.56.0.0/13 -j DROP
iptables -A INPUT -s 120.24.0.0/14 -j DROP
iptables -A INPUT -s 185.81.152.0/22 -j DROP
iptables -A INPUT -s 188.132.128.0/17 -j DROP
iptables -A INPUT -s 211.144.0.0/13 -j DROP
iptables -A INPUT -s 121.196.0.0/14 -j DROP
iptables -A INPUT -s 121.40.0.0/14 -j DROP
iptables -A INPUT -s 121.240.0.0/13 -j DROP
